# supervisor_default.ini
[program:php-fpm]
process_name = %(program_name)s
command = /usr/sbin/php-fpm7 -R --nodaemonize
autostart = true
autorestart = true
numprocs = 1
redirect_stderr = true

[program:nginx]
process_name = %(program_name)s
command = /usr/sbin/nginx -g 'daemon off;'
autostart = true
autorestart = true
numprocs = 1
redirect_stderr = true
